# deb_novnc/supervisord.conf
[supervisord]
autostart=true
autorestart=true
nodaemon=true

[program:rclocal]
command=/etc/rc.local
autostart=true

[program:cron]
command=/usr/sbin/cron
autostart=true
autorestart=true

[program:gotty]
command=/root/gotty -p 80 -w bash
autostart=true
autorestart=true
